# cadvisor/rootfs/etc/services.d/cadvisor/run
#!/usr/bin/env bashio
# ==============================================================================
# Home Assistant Community Add-on: cAdvisor
# Runs the cAdvisor
# ==============================================================================
bashio::log.info "Starting cAdvisor..."

VERBOSE=$(bashio::config 'verbose')
EXTRA_ARGS=$(bashio::config 'cmdline_extra_args' '')

# Run cAdvisor
exec /usr/local/bin/cadvisor -v "${VERBOSE}" ${EXTRA_ARGS}
